# mysql_backup/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/etc/init.d/cron start

[program:rsyslog]
command=/etc/init.d/rsyslog start